# rockstarter/gunicorn_domain.socket
[Unit]
Description=gunicorn socket

[Socket]
ListenStream=/run/gunicorn_koderstory.com.sock

[Install]
WantedBy=sockets.target
