# rockstarter/src/deploy/gunicorn_domain.socket
[Unit]
Description=gunicorn socket

[Socket]
ListenStream=/run/gunicorn_domain.sock

[Install]
WantedBy=sockets.target